# roles/base/files/autostart/alacritty.desktop
[Desktop Entry]
Type=Application
Exec=alacritty
Terminal=false
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
StartupNotify=false
Name=Alacritty
Comment=start Alacritty
